# supervisord-wine.conf
[program:explorer]
command=/opt/wine-stable/bin/wine /opt/wine-stable/lib/wine/i386-windows/explorer.exe
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
